# data/etc/init/datacollector.conf
description "Lacework datacollector"

start on runlevel [2345]
stop on runlevel [!2345]
limit nofile 524288 1048576
limit nproc 524288 1048576

respawn
post-stop exec sleep 10

kill timeout 20

script
	# modify these in /etc/default/$UPSTART_JOB (/etc/default/datacollector)
	LWDC=/var/lib/lacework/$UPSTART_JOB
	LWDC_OPTS=
	if [ -f /etc/default/$UPSTART_JOB ]; then
		. /etc/default/$UPSTART_JOB
	fi
	exec "$LWDC" $LWDC_OPTS >> /var/log/lacework/${UPSTART_JOB}.console.log 2>&1
end script

post-start script
	LWDC_OPTS=
	if [ -f /etc/default/$UPSTART_JOB ]; then
		. /etc/default/$UPSTART_JOB
	fi
end script
